# 02-playbooks/02-install-mattermost/install_mattermost/roles/web/templates/mattermost.upstart
start on runlevel [2345]
stop on runlevel [016]
respawn
limit nofile 50000 50000
chdir {{ mattermost_directory }}
setuid {{ mattermost_sysuser }}
exec bin/platform
